# daily-7am-claude/com.theodor349.daily7am.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.theodor349.daily7am</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/theodor349/.claude/daily-7am-claude/run.sh</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/theodor349/.claude/daily-7am-claude/logs/launchd.out</string>
    <key>StandardErrorPath</key>
    <string>/Users/theodor349/.claude/daily-7am-claude/logs/launchd.err</string>
</dict>
</plist>
